# bin/chard_baguette.conf
description "Start Baguette on boot"
# debug: grep chard_baguette /var/log/messages | tail -10
start on started ui
stop on stopping ui
task
script
    logger -t chard_baguette "starting chard_baguette"
    DEBUG=0
    [ -f "/home/chronos/user/.chard_baguette_debug" ] && DEBUG=1
    KEEP_RUNNING=0
    [ -f "/home/chronos/user/.chard_baguette_keep_running" ] && KEEP_RUNNING=1
    BASHRC=""
    attempt=0
    max_attempts=10
    wait_seconds=60
    while [ "$attempt" -lt "$max_attempts" ]; do
        logger -t chard_baguette "waiting for .bashrc (attempt $((attempt+1))/$max_attempts)"
        i=0
        while [ "$i" -lt "$wait_seconds" ]; do
            if [ -f "/home/chronos/user/.bashrc" ]; then
                BASHRC="/home/chronos/user/.bashrc"
                break
            fi
            sleep 1
            i=$((i + 1))
        done
        [ -n "$BASHRC" ] && break
        attempt=$((attempt + 1))
        sleep 5
    done
    if [ -z "$BASHRC" ]; then
        logger -t chard_baguette "timed out waiting for .bashrc"
        exit 1
    fi
    attempt=0
    max_attempts=10
    started=0
    while [ "$attempt" -lt "$max_attempts" ]; do
        logger -t chard_baguette "starting termina (attempt $((attempt+1))/$max_attempts)"
        if vmc start --enable-gpu --vm-type BAGUETTE --no-shell termina; then
            started=1
            break
        fi
        attempt=$((attempt + 1))
        logger -t chard_baguette "vmc start failed, retrying in 10s"
        sleep 10
    done
    if [ "$started" -ne 1 ]; then
        logger -t chard_baguette "failed to start termina after $max_attempts attempts"
        exit 1
    fi
    ready=0
    i=0
    while [ "$i" -lt 30 ]; do
        if timeout 10 vsh termina penguin -- true 2>/dev/null; then
            ready=1
            break
        fi
        sleep 2
        i=$((i + 1))
    done
    logger -t chard_baguette "guest reachable: $ready"
    if [ "$DEBUG" -eq 1 ]; then
        OUT="/home/chronos/user/MyFiles/chard_baguette_debug/$(date +%Y%m%d_%H%M%S)"
        mkdir -p "$OUT"
        vmc list > "$OUT/vmc_list.txt" 2>&1
        HASH=$(ls /run/daemon-store/crosvm/ 2>/dev/null | head -1)
        if [ -n "$HASH" ]; then
            cp "/run/daemon-store/crosvm/$HASH/log/dGVybWluYQ==.log" "$OUT/termina_crosvm.log" 2>/dev/null
        fi
        grep -E 'concierge|cicerone|crostini|garcon|CONTAINER_SETUP' /var/log/messages | tail -300 > "$OUT/host_messages.txt" 2>/dev/null
        vsh termina penguin -- env > "$OUT/guest_env.txt" 2>&1
        vsh termina penguin -- systemctl --user status garcon > "$OUT/guest_garcon_status.txt" 2>&1
        vsh termina penguin -- loginctl list-sessions > "$OUT/guest_sessions.txt" 2>&1
        logger -t chard_baguette "snapshot -> $OUT"
    fi
    if [ "$KEEP_RUNNING" -eq 0 ]; then
        logger -t chard_baguette "readying crostini"
        vmc stop termina
    else
        logger -t chard_baguette "run flag set"
    fi
    exit 0
end script
