# bin/chard_bridge_daemon.service
[Unit]
Description=Chard Bridge Daemon
After=multi-user.target
[Service]
Type=simple
ExecStart=/bin/chard_bridge_daemon
[Install]
WantedBy=multi-user.target
